# unified-engine/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:aws-engine]
command=uvicorn aws_api_server:app --host 0.0.0.0 --port 8000 --timeout-keep-alive 75
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/aws-engine.err.log
stdout_logfile=/var/log/supervisor/aws-engine.out.log
environment=PYTHONPATH="/app"

[program:compliance-engine]
command=uvicorn compliance_engine.api_server:app --host 0.0.0.0 --port 8001
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/compliance-engine.err.log
stdout_logfile=/var/log/supervisor/compliance-engine.out.log
environment=PYTHONPATH="/app"

[program:rule-engine]
command=uvicorn rule_api_server:app --host 0.0.0.0 --port 8002
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/rule-engine.err.log
stdout_logfile=/var/log/supervisor/rule-engine.out.log
environment=PYTHONPATH="/app"

